# etc/rc.d/bgpd
#!/bin/sh
#
# $OpenBSD$

daemon="/usr/sbin/bgpd"

. /etc/rc.d/rc.subr

pexp="bgpd: parent.*"

rc_cmd $1
