# etc/rc.d/dhcrelay
#!/bin/sh
#
# $OpenBSD$

daemon="/usr/sbin/dhcrelay"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
