# etc/rc.d/ftpd
#!/bin/sh
#
# $OpenBSD$

daemon="/usr/libexec/ftpd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
