# etc/rc.d/ftpproxy6
#!/bin/sh
#
# $OpenBSD$

daemon="/usr/sbin/ftp-proxy -6"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
