# etc/rc.d/httpd
#!/bin/sh
#
# $OpenBSD$

daemon="/usr/sbin/httpd"

. /etc/rc.d/rc.subr

pexp="httpd: parent.*"

rc_cmd $1
