# etc/rc.d/ldapd
#!/bin/sh
#
# $OpenBSD$

daemon="/usr/sbin/ldapd"

. /etc/rc.d/rc.subr

rc_cmd $1
