# etc/rc.d/ldomd
#!/bin/sh
#
# $OpenBSD$

daemon="/usr/sbin/ldomd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
