# etc/rc.d/ldpd
#!/bin/sh
#
# $OpenBSD$

daemon="/usr/sbin/ldpd"

. /etc/rc.d/rc.subr

pexp="ldpd: parent.*"

rc_cmd $1
