# etc/rc.d/lpd
#!/bin/sh
#
# $OpenBSD$

daemon="/usr/sbin/lpd"

. /etc/rc.d/rc.subr

rc_cmd $1
