# etc/rc.d/mopd
#!/bin/sh
#
# $OpenBSD$

daemon="/usr/sbin/mopd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_pre() {
	[ -d /tftpboot/mop ] || return 1
}

rc_cmd $1
