# etc/rc.d/mountd
#!/bin/sh
#
# $OpenBSD$

daemon="/sbin/mountd"

. /etc/rc.d/rc.subr

rc_cmd $1
