# etc/rc.d/named
#!/bin/sh
#
# $OpenBSD$

daemon="/usr/sbin/named"

. /etc/rc.d/rc.subr

pexp="named: \[priv\]"

rc_cmd $1
