# etc/rc.d/npppd
#!/bin/sh
#
# $OpenBSD$

daemon="/usr/sbin/npppd"

. /etc/rc.d/rc.subr

pexp="npppd: main"

rc_cmd $1
