# etc/rc.d/ntpd
#!/bin/sh
#
# $OpenBSD$

daemon="/usr/sbin/ntpd"

. /etc/rc.d/rc.subr

pexp="ntpd: \[priv\]"

rc_cmd $1
