# etc/rc.d/pflogd
#!/bin/sh
#
# $OpenBSD$

daemon="/sbin/pflogd"

. /etc/rc.d/rc.subr

pexp="pflogd: \[priv\]"

rc_cmd $1
