# etc/rc.d/portmap
#!/bin/sh
#
# $OpenBSD$

daemon="/usr/sbin/portmap"

. /etc/rc.d/rc.subr

rc_reload=NO
rc_restart=NO

rc_cmd $1
