# etc/rc.d/rarpd
#!/bin/sh
#
# $OpenBSD$

daemon="/usr/sbin/rarpd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_pre() {
	[ -s /etc/ethers ] || return 1
}

rc_cmd $1
