# etc/rc.d/relayd
#!/bin/sh
#
# $OpenBSD$

daemon="/usr/sbin/relayd"

. /etc/rc.d/rc.subr

pexp="relayd: parent.*"

rc_cmd $1
