# etc/rc.d/route6d
#!/bin/sh
#
# $OpenBSD$

daemon="/usr/sbin/route6d"

. /etc/rc.d/rc.subr

rc_cmd $1
