# etc/rc.d/sasyncd
#!/bin/sh
#
# $OpenBSD$

daemon="/usr/sbin/sasyncd"

. /etc/rc.d/rc.subr

pexp="sasyncd: \[priv\]"

rc_cmd $1
