# etc/rc.d/sndiod
#!/bin/sh
#
# $OpenBSD$

daemon="/usr/bin/sndiod"

. /etc/rc.d/rc.subr

rc_cmd $1
