# etc/rc.d/snmpd
#!/bin/sh
#
# $OpenBSD$

daemon="/usr/sbin/snmpd"

. /etc/rc.d/rc.subr

pexp="snmpd: parent.*"

rc_cmd $1
