# etc/rc.d/sshd
#!/bin/sh
#
# $OpenBSD$

daemon="/usr/sbin/sshd"

. /etc/rc.d/rc.subr

rc_cmd $1
