# etc/rc.d/statd
#!/bin/sh
#
# $OpenBSD$

daemon="/usr/sbin/rpc.statd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
