# etc/rc.d/syslogd
#!/bin/sh
#
# $OpenBSD: syslogd,v 1.3 2014/03/24 09:24:09 sthen Exp $

daemon="/usr/sbin/syslogd"

. /etc/rc.d/rc.subr

pexp="syslogd: \[priv\]"

rc_pre() {
	rm -f /dev/log
}

rc_cmd $1
