# etc/rc.d/unwind
#!/bin/ksh
#
# $OpenBSD$

daemon="/sbin/unwind"

. /etc/rc.d/rc.subr

rc_pre() {
	/usr/sbin/unbound-anchor -a /etc/unwind/trustanchor/root.key || true
}

rc_cmd $1
