# etc/rc.d/xdm
#!/bin/sh
#
# $OpenBSD$

daemon="/usr/X11R6/bin/xdm"

. /etc/rc.d/rc.subr

rc_cmd $1
